# vless-instance-b.service
# VLESS Instance B - Blue-Green Deployment
# Location: /etc/systemd/system/vless-instance-b.service
# This service runs on port 10081 and can be active or standby

[Unit]
Description=VLESS Reality Instance B (Port 10081)
Documentation=https://github.com/XTLS/Xray-core
After=network-online.target nss-lookup.target
Wants=network-online.target

[Service]
Type=simple
User=vless
Group=vless

# Binary and config paths (adjust to your installation)
ExecStart=/usr/local/bin/xray run -c /etc/vless/config-b.json

# Restart policy - critical for high availability
Restart=on-failure
RestartSec=5s
StartLimitBurst=5
StartLimitIntervalSec=60s

# Resource Limits - prevent one instance from consuming all resources
# Memory: Hard limit 512MB, soft warning at 400MB
MemoryMax=512M
MemoryHigh=400M
MemoryAccounting=yes

# CPU: Weight-based priority (default is 100, we use 200 for higher priority)
CPUWeight=200
CPUAccounting=yes

# Nice value: -20 to 19, lower = higher priority (-5 is elevated)
Nice=-5

# File descriptor limit - important for high connection counts
LimitNOFILE=65535

# Process limits
LimitNPROC=512

# Core dumps disabled for security
LimitCORE=0

# Security hardening
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/vless

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vless-instance-b

# Environment
Environment="XRAY_LOCATION_ASSET=/usr/local/share/xray"

[Install]
WantedBy=multi-user.target
